# clickhouse/systemd/clickhouse-server.service
[Unit]
Description=ClickHouse Server (custom binary)
After=network.target

[Service]
Type=simple
User=clickhouse
Group=clickhouse
ExecStart=/usr/local/bin/clickhouse-server --config-file=/etc/clickhouse-server/config.xml --log-file=/var/log/clickhouse/server.log --errorlog-file=/var/log/clickhouse/server.err.log
Restart=on-failure
RestartSec=2s
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
